# crontab
0 */4 * * * /usr/bin/python3 /home/pi/price_monitor/update.py  >> /home/pi/price_monitor/monitor.log 2>&1
